# run.conf
[program:daemon]
command=/neutrino/run.sh > /neutrino/logs/startup.log 2>&1
